# res/debian/prerm.in
#!/bin/sh
set +e
systemctl stop @PROJECT_NAME@
systemctl disable @PROJECT_NAME@
